# prep/debian/postrm
#!/bin/sh
#
# postrm script for gossimon
set -e

echo "Running gossimon postrm 1:[$1] 2:[$2]"

if [ "$1" = purge  -o "$1" = "remove" ]; then
    update-rc.d infod remove
fi
